# experimental/noah/noah.service
[Unit]
Description=Noah daily scripture ingestion
After=network.target

[Service]
Type=oneshot
User=%i
WorkingDirectory=/home/ark/orthodoxphronema
ExecStart=/usr/bin/python3 experimental/noah/harness.py
Environment=NOAH_VAULT_DIR=/home/ark/noah-vault
StandardOutput=journal
StandardError=journal
TimeoutStartSec=900
